# configs/services/fne-router.service
[Unit]
Description=DVMFNE Router
After=network.target

[Service]
RestartSec=2s
Type=simple
WorkingDirectory=/usr/local/share/dvmfne
ExecStart=/usr/bin/python3 /usr/local/share/dvmfne/dvmfne/fne_router.py -c /etc/dvmfne/fne_router.cfg
Restart=always

[Install]
WantedBy=multi-user.target
